# refactored/DataServer/data/deploy/velox-server.service
# Velox Master Server - Systemd Service
# =============================================
# Install: sudo ./data/deploy/install-server.sh
# Manual:  sudo cp data/deploy/velox-server.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now velox-server
#
# Logs:    journalctl -u velox-server -f
# Status:  systemctl status velox-server
#
# Environment file: /etc/velox-server.env
#   Configure all VELOX_* vars there (see data/deploy/velox-server.env for template).
# =============================================

[Unit]
Description=Velox Master Server
Documentation=https://codebuff.com/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# User and group to run as (create with: sudo useradd -r -s /usr/sbin/nologin -m velox)
User=velox
Group=velox

# Working directory (project root containing DataServer/, RemoteCodex/, etc.)
WorkingDirectory=/opt/velox/current

# Environment file with all VELOX_* config vars
EnvironmentFile=/etc/velox-server.env

# The binary path (relative to WorkingDirectory)
ExecStart=/opt/velox/current/DataServer/bin/velox-server

# Restart behavior
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Graceful shutdown — Stop() tears down gRPC immediately (workers hold
# open bidirectional streams that would block GracefulStop forever).
# 60s is a generous safety net; the actual shutdown completes in <5s.
TimeoutStopSec=60
KillMode=process
KillSignal=SIGTERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Sandboxing (hardening)
ProtectSystem=full
ProtectHome=false
NoNewPrivileges=true
PrivateTmp=true
CapabilityBoundingSet=
AmbientCapabilities=

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=velox-server

[Install]
WantedBy=multi-user.target
