# deploy/velox-server.service
# Velox Master Server - Systemd Service
# =============================================
# Install: sudo ./deploy/install-server.sh
# Manual:  sudo cp deploy/velox-server.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now velox-server
#
# Logs:    journalctl -u velox-server -f
# Status:  systemctl status velox-server
#
# Environment file: /etc/velox-server.env
#   Configure all VELOX_* vars there, including VELOX_SERVER_IMAGE pinned to
#   ghcr.io/<owner>/velox-server@sha256:<digest>.
# =============================================

[Unit]
Description=Velox Master Server
Documentation=https://codebuff.com/docs
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
# The unit drives Docker itself, so it must run as root. The application
# process still runs inside the container as the image's unprivileged user.
EnvironmentFile=/etc/velox-server.env
ExecStartPre=-/usr/bin/docker rm -f velox-server
ExecStartPre=/usr/bin/docker pull ${VELOX_SERVER_IMAGE}
ExecStart=/usr/bin/docker run --rm \
    --name velox-server \
    --network host \
    --env-file /etc/velox-server.env \
    --env VELOX_HTTP_PORT=${VELOX_MASTER_PORT} \
    --mount type=bind,src=/var/lib/velox,dst=/var/lib/velox \
    --mount type=bind,src=/etc/velox,dst=/etc/velox,readonly \
    ${VELOX_SERVER_IMAGE}
ExecStop=/usr/bin/docker stop -t 60 velox-server
ExecStopPost=-/usr/bin/docker rm -f velox-server

# Restart behavior
Restart=always
RestartSec=10

# Graceful shutdown — Stop() tears down gRPC immediately (workers hold
# open bidirectional streams that would block GracefulStop forever).
# 60s is a generous safety net; the actual shutdown completes in <5s.
TimeoutStopSec=60
KillMode=process
KillSignal=SIGTERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=velox-server

[Install]
WantedBy=multi-user.target
